# overlays/systemd-repart/etc/systemd/systemd-repart.service.d/bugfix-25849.conf
[Unit]
Wants=modprobe@loop.service modprobe@dm_mod.service
After=initrd-usr-fs.target modprobe@loop.service modprobe@dm_mod.service